# systemd-userlevel/github.timer
[Unit]
Description=Timer for periodically removing github followers at user level

[Timer]
OnUnitActiveSec=60
